# docker-laravel/local/php/supervisor/vite-watchdog.conf
; Self-heal the Vite dev server (local only). Editing a watched Blade can crash
; vite (exit 1), deleting public/hot → ViteManifestNotFoundException 500s. Poll
; every 5s; if npm-dev is RUNNING but public/hot has been gone for ~10s (2-strike
; grace — vite boots in <1s, so 10s of RUNNING-without-hot means a crash, not a
; normal startup), restart npm-dev. Recovery window ≈10s.
[program:vite-watchdog]
command=/bin/sh -c 'm=0; while sleep 5; do if supervisorctl -c /etc/supervisor/supervisord.conf status npm-dev | grep -q RUNNING; then if [ ! -f /var/www/public/hot ]; then m=$((m+1)); else m=0; fi; if [ "$m" -ge 2 ]; then echo "[vite-watchdog] public/hot missing — restarting npm-dev"; supervisorctl -c /etc/supervisor/supervisord.conf restart npm-dev; m=0; fi; else m=0; fi; done'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
